# debian/sillybear.conffiles
/etc/init.d/sillybear
/etc/sillybear/run
/etc/sillybear/log/run
