# debian/sillybear.prerm
#!/bin/sh
set -e

test "$1" = 'remove' || test "$1" = 'deconfigure' || exit 0
if test -x /etc/init.d/sillybear; then
  if test -x /usr/sbin/invoke-rc.d; then
    invoke-rc.d sillybear stop
  else
    /etc/init.d/sillybear stop
  fi
fi
